# deploy/systemd/agent-server-cosyvoice-fastapi.service
[Unit]
Description=agent-server CosyVoice FastAPI TTS
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=ubuntu
WorkingDirectory=/home/ubuntu/agent-server
EnvironmentFile=-/etc/agent-server/cosyvoice-fastapi.env
ExecStart=/home/ubuntu/agent-server/scripts/run-cosyvoice-fastapi-local.sh
Restart=always
RestartSec=3

[Install]
WantedBy=multi-user.target
